# plugins/os-deepinspector/src/etc/rc.d/deepinspector
#!/bin/sh
# deepinspector - Deep Packet Inspector daemon script

. /etc/rc.subr

name="deepinspector"
rcvar=deepinspector_enable
pidfile="/var/run/${name}.pid"
command="/usr/local/bin/python3.11"
command_args="/usr/local/opnsense/scripts/OPNsense/DeepInspector/deepinspector_engine.py"
procname="/usr/local/bin/python3.11"

load_rc_config $name

: ${deepinspector_enable:=NO}

start_precmd="${name}_prestart"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"

deepinspector_prestart()
{
    # Create necessary directories
    [ -d /var/run ] || install -d -m 755 /var/run
    [ -d /var/log/deepinspector ] || install -d -m 755 /var/log/deepinspector
    [ -d /usr/local/etc/deepinspector ] || install -d -m 755 /usr/local/etc/deepinspector
    
    # Check if Python3 exists
    if [ ! -f "$command" ]; then
        echo "Error: Python3 not found at $command"
        return 1
    fi
    
    # Check if engine script exists
    if [ ! -f "$command_args" ]; then
        echo "Error: Engine script not found at $command_args"
        return 1
    fi
    
    return 0
}

deepinspector_start()
{
    echo "[$(date)] Starting ${name}..." >> /var/log/deepinspector/daemon.log
    
    # Check if already running
    if [ -f "${pidfile}" ]; then
        OLD_PID=$(cat ${pidfile})
        if kill -0 "${OLD_PID}" 2>/dev/null; then
            echo "[$(date)] ${name} is already running with PID ${OLD_PID}." >> /var/log/deepinspector/daemon.log
            return 1
        else
            echo "[$(date)] Removing stale PID file" >> /var/log/deepinspector/daemon.log
            rm -f "${pidfile}"
        fi
    fi
    
    # Export configuration from OPNsense
    if [ -f /usr/local/opnsense/scripts/OPNsense/DeepInspector/export_config.py ]; then
        echo "[$(date)] Exporting configuration..." >> /var/log/deepinspector/daemon.log
        /usr/local/bin/python3.11 /usr/local/opnsense/scripts/OPNsense/DeepInspector/export_config.py >> /var/log/deepinspector/daemon.log 2>&1
    else
        echo "[$(date)] Warning: export_config.py not found, using default configuration" >> /var/log/deepinspector/daemon.log
    fi
    
    # Start the process using daemon utility for better process management
    echo "[$(date)] Starting engine process..." >> /var/log/deepinspector/daemon.log
    
    # Use daemon utility if available, otherwise fallback to manual background process
    if command -v daemon >/dev/null 2>&1; then
        /usr/sbin/daemon -p "${pidfile}" -f "${command}" "${command_args}"
    else
        # Fallback method
        "${command}" "${command_args}" > /var/log/deepinspector/stdout.log 2> /var/log/deepinspector/stderr.log &
        echo $! > "${pidfile}"
    fi
    
    sleep 3
    
    # Check if process is still running
    if [ -f "${pidfile}" ]; then
        NEW_PID=$(cat ${pidfile})
        if kill -0 "${NEW_PID}" 2>/dev/null; then
            echo "[$(date)] ${name} started successfully with PID ${NEW_PID}" >> /var/log/deepinspector/daemon.log
            return 0
        else
            echo "[$(date)] Failed to start ${name} (PID ${NEW_PID} not alive)" >> /var/log/deepinspector/daemon.log
            rm -f "${pidfile}"
            return 1
        fi
    else
        echo "[$(date)] Failed to start ${name} (no PID file created)" >> /var/log/deepinspector/daemon.log
        return 1
    fi
}

deepinspector_stop()
{
    echo "[$(date)] Stopping ${name}..." >> /var/log/deepinspector/daemon.log
    
    if [ -f "${pidfile}" ]; then
        PID=$(cat "${pidfile}")
        if kill -0 "${PID}" 2>/dev/null; then
            # Send TERM signal for graceful shutdown
            echo "[$(date)] Sending TERM signal to PID ${PID}" >> /var/log/deepinspector/daemon.log
            kill -TERM "${PID}"
            
            # Wait for graceful shutdown with timeout
            TIMEOUT=10
            while [ $TIMEOUT -gt 0 ] && kill -0 "${PID}" 2>/dev/null; do
                sleep 1
                TIMEOUT=$((TIMEOUT - 1))
            done
            
            # Force kill if still running
            if kill -0 "${PID}" 2>/dev/null; then
                echo "[$(date)] Process still running, sending KILL signal" >> /var/log/deepinspector/daemon.log
                kill -KILL "${PID}" 2>/dev/null
                sleep 2
            fi
            
            echo "[$(date)] ${name} stopped." >> /var/log/deepinspector/daemon.log
            rm -f "${pidfile}"
            return 0
        else
            echo "[$(date)] PID ${PID} not running, removing stale PID file" >> /var/log/deepinspector/daemon.log
            rm -f "${pidfile}"
            return 0
        fi
    else
        echo "[$(date)] No PID file found. Is ${name} running?" >> /var/log/deepinspector/daemon.log
        return 1
    fi
}

deepinspector_status()
{
    if [ -f "${pidfile}" ]; then
        PID=$(cat "${pidfile}")
        if kill -0 "${PID}" 2>/dev/null; then
            echo "${name} is running as PID ${PID}"
            
            # Get additional status info
            if [ -f /var/log/deepinspector/stats.json ]; then
                echo "Statistics available in /var/log/deepinspector/stats.json"
                # Show basic stats if jq is available
                if command -v jq >/dev/null 2>&1; then
                    PACKETS=$(jq -r '.packets_analyzed // 0' /var/log/deepinspector/stats.json 2>/dev/null)
                    THREATS=$(jq -r '.threats_detected // 0' /var/log/deepinspector/stats.json 2>/dev/null)
                    echo "Packets analyzed: ${PACKETS}, Threats detected: ${THREATS}"
                fi
            fi
            
            # Show memory usage
            MEMORY=$(ps -o rss= -p "${PID}" 2>/dev/null | awk '{print int($1/1024)"MB"}')
            if [ -n "$MEMORY" ]; then
                echo "Memory usage: ${MEMORY}"
            fi
            
            return 0
        else
            echo "${name} PID file exists but process is not running"
            return 1
        fi
    else
        echo "${name} is not running"
        return 1
    fi
}

deepinspector_restart()
{
    echo "[$(date)] Restarting ${name}..." >> /var/log/deepinspector/daemon.log
    deepinspector_stop
    sleep 3
    deepinspector_start
}

# Handle different invocation methods
case "$1" in
    start)
        deepinspector_prestart && deepinspector_start
        ;;
    stop)
        deepinspector_stop
        ;;
    restart)
        deepinspector_restart
        ;;
    status)
        deepinspector_status
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac