# plugins/os-siemlogger/src/etc/rc.d/siemlogger
#!/bin/sh
# siemlogger - SiemLogger daemon script

# Carica le funzioni di rc.subr per la gestione dei daemon
. /etc/rc.subr

name="siemlogger"
rcvar="${name}_enable"
pidfile="/var/run/${name}.pid"
command="/usr/local/bin/python3.11"
command_args="/usr/local/opnsense/scripts/OPNsense/SiemLogger/siemlogger_engine.py"
config_file="/usr/local/etc/siemlogger/config.json"
log_dir="/var/log/siemlogger"
log_file="${log_dir}/service.log"

# Impostazioni predefinite
: ${siemlogger_enable:=NO}

# Funzioni specifiche
siemlogger_prestart()
{
    # Crea directory necessarie
    [ -d /var/run ] || install -d -m 755 /var/run 2>/dev/null
    [ -d "${log_dir}" ] || install -d -m 750 "${log_dir}" 2>/dev/null
    [ -d /usr/local/etc/siemlogger ] || install -d -m 755 /usr/local/etc/siemlogger 2>/dev/null
    [ -d /var/db/siemlogger ] || install -d -m 755 /var/db/siemlogger 2>/dev/null

    # Verifica Python
    if [ ! -x "${command}" ]; then
        echo "[$(date '+%Y-%m-%d %H:%M:%S')] Error: Python3 not found or not executable at ${command}" >> "${log_file}" 2>/dev/null
        return 1
    fi

    # Verifica engine script
    if [ ! -x "${command_args}" ]; then
        echo "[$(date '+%Y-%m-%d %H:%M:%S')] Error: Engine script not found or not executable at ${command_args}" >> "${log_file}" 2>/dev/null
        return 1
    fi

    # Verifica moduli
    REQUIRED_MODULES="settings_logger export_events health_check"
    SCRIPTS_DIR="/usr/local/opnsense/scripts/OPNsense/SiemLogger"
    for module in ${REQUIRED_MODULES}; do
        if [ ! -f "${SCRIPTS_DIR}/${module}.py" ]; then
            echo "[$(date '+%Y-%m-%d %H:%M:%S')] Error: Required module not found: ${SCRIPTS_DIR}/${module}.py" >> "${log_file}" 2>/dev/null
            return 1
        fi
        chmod +x "${SCRIPTS_DIR}/${module}.py" 2>/dev/null || true
    done

    # Crea config se assente
    if [ ! -f "${config_file}" ]; then
        echo "[$(date '+%Y-%m-%d %H:%M:%S')] Creating basic config file..." >> "${log_file}" 2>/dev/null
        cat > "${config_file}" << 'EOF'
{
  "general": {"enabled": true, "log_level": "INFO"},
  "siem_export": {"export_enabled": false, "export_format": "json"}
}
EOF
        chmod 644 "${config_file}" 2>/dev/null
    fi

    return 0
}

siemlogger_start()
{
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] Starting ${name}..." >> "${log_file}" 2>/dev/null

    if [ -f "${pidfile}" ]; then
        OLD_PID=$(cat "${pidfile}")
        if kill -0 "${OLD_PID}" 2>/dev/null; then
            echo "[$(date '+%Y-%m-%d %H:%M:%S')] ${name} is already running with PID ${OLD_PID}." >> "${log_file}" 2>/dev/null
            echo "${name} is already running as PID ${OLD_PID}"
            return 1
        else
            echo "[$(date '+%Y-%m-%d %H:%M:%S')] Removing stale PID file" >> "${log_file}" 2>/dev/null
            rm -f "${pidfile}" 2>/dev/null
        fi
    fi

    # Avvia il processo
    if command -v daemon >/dev/null 2>&1; then
        /usr/sbin/daemon -p "${pidfile}" -f "${command}" "${command_args}" >> "${log_file}" 2>&1
    else
        "${command}" "${command_args}" > "${log_dir}/stdout.log" 2> "${log_dir}/stderr.log" &
        echo $! > "${pidfile}" 2>/dev/null
    fi

    sleep 2

    if [ -f "${pidfile}" ]; then
        NEW_PID=$(cat "${pidfile}")
        if kill -0 "${NEW_PID}" 2>/dev/null; then
            echo "[$(date '+%Y-%m-%d %H:%M:%S')] ${name} started successfully with PID ${NEW_PID}" >> "${log_file}" 2>/dev/null
            echo "${name} started successfully as PID ${NEW_PID}"
            return 0
        else
            echo "[$(date '+%Y-%m-%d %H:%M:%S')] Failed to start ${name} (PID ${NEW_PID} not alive)" >> "${log_file}" 2>/dev/null
            echo "[$(date '+%Y-%m-%d %H:%M:%S')] Check ${log_dir}/stderr.log for errors" >> "${log_file}" 2>/dev/null
            rm -f "${pidfile}" 2>/dev/null
            return 1
        fi
    else
        echo "[$(date '+%Y-%m-%d %H:%M:%S')] Failed to start ${name} (no PID file created)" >> "${log_file}" 2>/dev/null
        return 1
    fi
}

siemlogger_stop()
{
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] Stopping ${name}..." >> "${log_file}" 2>/dev/null

    if [ -f "${pidfile}" ]; then
        PID=$(cat "${pidfile}")
        if kill -0 "${PID}" 2>/dev/null; then
            kill -TERM "${PID}" 2>/dev/null
            sleep 5
            if kill -0 "${PID}" 2>/dev/null; then
                kill -KILL "${PID}" 2>/dev/null
            fi
            rm -f "${pidfile}" 2>/dev/null
            echo "[$(date '+%Y-%m-%d %H:%M:%S')] ${name} stopped." >> "${log_file}" 2>/dev/null
            echo "${name} stopped"
            return 0
        else
            rm -f "${pidfile}" 2>/dev/null
            echo "[$(date '+%Y-%m-%d %H:%M:%S')] PID ${PID} not running, removed stale PID file" >> "${log_file}" 2>/dev/null
            return 0
        fi
    else
        echo "[$(date '+%Y-%m-%d %H:%M:%S')] No PID file found. Is ${name} running?" >> "${log_file}" 2>/dev/null
        echo "${name} is not running"
        return 1
    fi
}

siemlogger_status()
{
    if [ -f "${pidfile}" ]; then
        PID=$(cat "${pidfile}")
        if kill -0 "${PID}" 2>/dev/null; then
            echo "${name} is running as PID ${PID}"
            return 0
        else
            echo "${name} PID file exists but process is not running"
            return 1
        fi
    else
        echo "${name} is not running"
        return 1
    fi
}

siemlogger_restart()
{
    siemlogger_stop
    sleep 1
    siemlogger_start
}

# Definisci i comandi
start_precmd="${name}_prestart"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"

run_rc_command "$1"